# deployment/remote-ws-module/update-gcp-dns.service
[Unit]
Description=Update DNS of instance
After=network.target

[Service]
Restart=on-failure
RestartSec=5s
EnvironmentFile=-/var/update-gcp-dns.env
ExecStart=/usr/local/bin/update-gcp-dns
StandardOutput=journal

[Install]
WantedBy=multi-user.target
